# data/openfortivpn-gui-helper.service
[Unit]
Description=OpenFortiVPN GUI Helper Daemon
Documentation=https://github.com/shini4i/openfortivpn-gui
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/openfortivpn-gui-helper
Restart=on-failure
RestartSec=5
Group=openfortivpn-gui

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Restrict Linux capabilities to the minimum required for VPN operations
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_READ_SEARCH

# Runtime and state directories
RuntimeDirectory=openfortivpn-gui
RuntimeDirectoryMode=0750
StateDirectory=openfortivpn-gui
StateDirectoryMode=0750

# Watchdog
WatchdogSec=30

[Install]
WantedBy=multi-user.target
