# data/openfortivpn-gui-helper.service
[Unit]
Description=OpenFortiVPN GUI Helper Daemon
Documentation=https://github.com/shini4i/openfortivpn-gui
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/openfortivpn-gui-helper
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Runtime and state directories
RuntimeDirectory=openfortivpn-gui
RuntimeDirectoryMode=0755
StateDirectory=openfortivpn-gui
StateDirectoryMode=0755

# Watchdog
WatchdogSec=30

[Install]
WantedBy=multi-user.target
